# server-config/fileserver.service
[Unit]
Description=BuildCat File Server — Flask REST API for file management
After=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/var/www/chat
ExecStart=/usr/bin/python3 /var/www/chat/fileserver.py
# 环境变量（Gateway Token 从 openclaw.json 读取，此处可覆盖）
# Environment=FILESERVER_TOKEN=your-token-here
Environment=PORT=5050
Environment=FLASK_DEBUG=0

# 自动重启策略：进程 crash（非零退出码、信号终止）后自动重启
# RestartSec 重启前等待秒数，避免快速循环重启
Restart=on-failure
RestartSec=5
# 可选：启用 systemd 主动健康检查（需应用支持 sd_notify）
# WatchdogSec=30

# 退出后等待时间
TimeoutStopSec=10

# 日志输出到 journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=fileserver

# 安全加固
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/ubuntu/.openclaw/workspace-assistant
ReadWritePaths=/home/ubuntu/.openclaw/workspace-build-cat
ReadWritePaths=/home/ubuntu/.openclaw/user-sessions
ReadWritePaths=/home/ubuntu/.openclaw/user-files

[Install]
WantedBy=multi-user.target
